# collector/launchd/com.readinggo.collector-prewarm.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  ReadingGo 시드 선충전 배치 (spec seed-collector.md §3.2) — 매일 새벽 04:00 KST.
  인기 카탈로그책을 seed_queue 에 low 우선순위로 큐잉(실제 크롤은 상시 poller 가 한가할 때 소진).
  설치 전 __NODE__ 와 __COLLECTOR_DIR__ 치환(README 참조).
  설치:  cp 후 launchctl load -w ~/Library/LaunchAgents/com.readinggo.collector-prewarm.plist
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.readinggo.collector-prewarm</string>

  <key>ProgramArguments</key>
  <array>
    <string>__NODE__</string>
    <string>__COLLECTOR_DIR__/prewarm.mjs</string>
  </array>

  <key>WorkingDirectory</key>
  <string>__COLLECTOR_DIR__</string>

  <!-- 매일 04:00(맥미니 로컬 시각 = KST). 데몬과 같은 브라우저 자원이라 새벽 분리. -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key><integer>4</integer>
    <key>Minute</key><integer>0</integer>
  </dict>

  <key>RunAtLoad</key>
  <false/>

  <key>StandardOutPath</key>
  <string>__COLLECTOR_DIR__/logs/prewarm.out.log</string>
  <key>StandardErrorPath</key>
  <string>__COLLECTOR_DIR__/logs/prewarm.err.log</string>
</dict>
</plist>
